# scripts/com.winkel.nextjs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.winkel.nextjs</string>
  <key>ProgramArguments</key>
  <array>
    <string>NODE_PATH</string>
    <string>PROJECT_PATH/node_modules/next/dist/bin/next</string>
    <string>start</string>
    <string>-p</string>
    <string>3001</string>
  </array>
  <key>WorkingDirectory</key>
  <string>PROJECT_PATH</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
  <key>StandardOutPath</key>
  <string>PROJECT_PATH/logs/nextjs.log</string>
  <key>StandardErrorPath</key>
  <string>PROJECT_PATH/logs/nextjs.log</string>
</dict>
</plist>
